# docker/app/backend/crontab
# * * * * * cd /var/www/html && php artisan schedule:run >> /dev/null 2>&1
* * * * * su -c "cd /var/www/html && php artisan schedule:run >> /dev/null 2>&1" -s /bin/bash www-data
